# supervisor.conf
[supervisord]
nodaemon=true

[program:dockerui]
command=/app/dockerui
directory=/app

[program:nginx]
command=/usr/sbin/nginx
